# offline/build/docker/config/supervisor.conf
[supervisord]
nodaemon=true
user=www-data

[inet_http_server]
port=127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0

## Does not work in docker -- removed to prevent errors/log spam
# [program:cron]
# # must run in foreground w/ supervisor or it will think cron failed
# command=cron -f
# autostart=true
# autorestart=false
# stdout_logfile = /dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile = /dev/stderr
# stderr_logfile_maxbytes=0

## Replaced by horizon
# [program:porta-worker]
# process_name=%(program_name)s_%(process_num)02d
# command=/usr/local/bin/php /var/www/artisan queue:work --sleep=0 --timeout=0 --memory=512 --tries=3
# autostart=true
# autorestart=true
# user=www-data
# # assign multiple workers to a queue and process jobs concurrently
# numprocs=2
# redirect_stderr=true
# stdout_logfile=/var/www/logs/default_queue.log

[program:porta-task-worker]
process_name=%(program_name)s_%(process_num)02d
# run scheduled tasks every second
command=/usr/local/bin/php /var/www/artisan schedule:run-cronless --frequency=1
autostart=true
autorestart=true
user=www-data
# assign multiple workers
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/logs/default_tasks.log

[program:porta-horizon-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /var/www/artisan horizon
autostart=true
autorestart=true
user=www-data
# assign multiple workers
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/logs/horizon.log
stopwaitsecs=3600
